# crontab
# Run at 6 AM ET daily (11 AM UTC, adjust based on DST)
05 6 * * * cd /app && /usr/local/bin/node /app/redeem.js >> /var/log/cron.log 2>&1
# Note: ET is UTC-5 (EST) or UTC-4 (EDT). Adjust the hour as needed.
